# apt-daily.timer.d/apt-update.conf
# /etc/systemd/system/apt-daily.timer.d/apt-update.conf

[Timer]
OnCalendar=06:00
RandomizedDelaySec=0