# apt-daily.timer.d/override.conf
# /etc/systemd/system/apt-daily.timer.d/override.conf

[Timer]
OnCalendar= # To properly reset the default value
OnCalendar=*-*-27 06:00:00
RandomizedDelaySec=0
