# apt-daily.timer.d/override.conf
# /etc/systemd/system/apt-daily.timer.d/override.conf

[Timer]
# To properly reset the default value
OnCalendar=
OnCalendar=Sat *-*~07/1 07:00:00
RandomizedDelaySec=0
