# apt-daily-upgrade.timer.d/apt-upgrade.conf
# /etc/systemd/system/apt-daily-upgrade.timer.d/apt-upgrade.conf

[Timer]
OnCalendar=Sat 06:30
RandomizedDelaySec=0