# apt-daily-upgrade.timer.d/override.conf
# /etc/systemd/system/apt-daily-upgrade.timer.d/override.conf

[Timer]
# To properly reset the default value
OnCalendar=
OnCalendar=Sat *-*~07/1 06:30:00
RandomizedDelaySec=0
