# apt-daily-upgrade.timer.d/override.conf
# /etc/systemd/system/apt-daily-upgrade.timer.d/override.conf

[Timer]
OnCalendar= # To properly reset the default value
OnCalendar=*-*-27 06:30:00
RandomizedDelaySec=0
