# apt-daily-upgrade.timer.d/override.conf
# /etc/systemd/system/apt-daily-upgrade.timer.d/override.conf

[Timer]
# To properly reset the default value
OnCalendar=
OnCalendar=*-*-* 7:45
RandomizedDelaySec=0
